# workspace/setup/auto_restart.service
[Unit]
Description=AI Collaboration Platform - Auto Restart Service
After=network.target

[Service]
Type=simple
User=moco
WorkingDirectory=/home/moco/workspace
Environment=PYTHONPATH=/home/moco/workspace
Environment=LOG_LEVEL=INFO
ExecStart=/usr/bin/python3 /home/moco/workspace/services/auto_restart_service.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
Restart=always
RestartSec=10
RestartPreventExitStatus=255
StandardOutput=append:/home/moco/workspace/logs/auto_restart_systemd.log
StandardError=append:/home/moco/workspace/logs/auto_restart_systemd_error.log

# リソース制限
LimitNOFILE=65535
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
